# src/bird/service/bird6/run
#!/bin/sh

if [ ! -f /etc/network/bird4.conf ]
then
	install -m0666 /etc/bird/bird6.conf /etc/network/bird6.conf
fi

mkdir -p /run/bird
exec /usr/sbin/bird6 -c /etc/network/bird6.conf -f
