# src/bird/service/bird4/run
#!/bin/sh

if [ ! -f /etc/network/bird4.conf ]
then
	install -m0666 /etc/bird/bird.conf /etc/network/bird4.conf
fi

mkdir -p /run/bird
exec /usr/sbin/bird -c /etc/network/bird4.conf -f
