# src/etc/service/lldpd/run
#!/bin/sh

exec chpst -0 lldpd -d >/dev/null 2>&1
