# src/etc/service/configwatch/run
#!/bin/sh

exec chpst -012 /etc/dockrouter/configwatch.sh
